# deploy/systemd/omnisync.service
[Unit]
Description=OmniSync FastAPI Service
After=network.target

[Service]
Type=simple
User=omnisync
Group=omnisync
WorkingDirectory=/opt/omnisync
EnvironmentFile=/opt/omnisync/deploy/ENV
ExecStart=/opt/omnisync/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
